# systemd/echo_face.service
[Unit]
Description=Project Echo Face Renderer
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=echo
Group=echo
WorkingDirectory=/opt/echo-ai
EnvironmentFile=/opt/echo-ai/.env
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=SDL_VIDEODRIVER=dummy
Environment=PYGAME_HIDE_SUPPORT_PROMPT=1
Environment=ECHO_FACE_FULLSCREEN=0
ExecStart=/opt/echo-ai/.venv/bin/python /opt/echo-ai/echo_face.py
Restart=on-failure
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target

